# upstart-script/mongo-image-server.conf
description "go web server"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 1000 1
exec sudo -u www-data /opt/image-server/mongo-image-server >> /var/log/nginx/image-server.staging.log